# deploy/systemd/pmg-portal.service
[Unit]
Description=PMG Portal (Gunicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pmg-portal/src
EnvironmentFile=/opt/pmg-portal/.env
ExecStart=/opt/pmg-portal/src/.venv/bin/gunicorn pmg_portal.wsgi:application --bind ${APP_BIND} --workers 2 --access-logfile -
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
